# backend/deployment/neo-alexandria-worker.service
[Unit]
Description=Neo Alexandria Edge Worker
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/neo-alexandria-2.0/backend
Environment="PATH=/path/to/neo-alexandria-2.0/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/path/to/neo-alexandria-2.0/backend/.venv/bin/python3 worker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
